# deploy/com.example.claude-mobile.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  Template launchd agent. Do not edit the paths by hand — run deploy/install.sh,
  which copies this file and substitutes __NODE__, __REPO__, and __HOME__ with
  values from your machine.
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>local.claude-mobile</string>

  <key>ProgramArguments</key>
  <array>
    <string>__NODE__</string>
    <string>__REPO__/server/dist/index.js</string>
  </array>

  <key>WorkingDirectory</key>
  <string>__REPO__</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <key>HOME</key>
    <string>__HOME__</string>
  </dict>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>__HOME__/Library/Logs/claude-mobile.out.log</string>
  <key>StandardErrorPath</key>
  <string>__HOME__/Library/Logs/claude-mobile.err.log</string>

  <key>ThrottleInterval</key>
  <integer>5</integer>
</dict>
</plist>
